# deployment/trading-bot.service
[Unit]
Description=TufureOS Trading Bot
After=network.target

[Service]
Type=simple
User=%I
WorkingDirectory=/home/%i/trading_bot
Environment=PYTHONPATH=/home/%i/trading_bot
EnvironmentFile=/home/%i/trading_bot/.env
ExecStart=/usr/bin/python3 -m src.bot --config config.json
Restart=on-failure
RestartSec=60
StandardOutput=append:/home/%i/trading_bot/logs/bot.log
StandardError=append:/home/%i/trading_bot/logs/bot_error.log

[Install]
WantedBy=multi-user.target
